# qr3_13.conf
[program:qr3_13]
command=python3 /home/curator3/qr3_13/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_13/test.err.log
stdout_logfile=/home/curator3/qr3_13/test.err.log
